# vgaswitch.service
[Unit]
Description=vgaswitch

[Service]
ExecStart=/usr/share/vgaswitch/vgaswitch.sh

[Install]
WantedBy=multi-user.target
